# recipes-bsp/board-init/board-init/board-init.service
[Unit]
Description=Run miscellaneous board init scripts
After=local-fs.target sysinit.target
DefaultDependencies=no

[Service]
Type=oneshot
User=root
Nice=-20
ExecStart=/usr/sbin/board-init.sh
RemainAfterExit=yes

[Install]
WantedBy=basic.target
