# recipes-support/ce-testing-scripts/ce-testing-scripts/imx8mm-jaguar-sentai/ce-iperf3-test.service
[Unit]
Description=CE iperf3 testing

[Service]
User=root
ExecStart=/bin/bash /usr/bin/ce-iperf3-test.sh

[Install]
WantedBy=multi-user.target
