# recipes-support/ce-testing-scripts/ce-testing-scripts/imx8mm-jaguar-sentai/ce-led-test.service
[Unit]
Description=CE LED testing

[Service]
User=root
ExecStart=/bin/bash /usr/bin/ce-led-test.sh

[Install]
WantedBy=multi-user.target
