# recipes-support/ce-testing-scripts/ce-testing-scripts/imx8mm-jaguar-sentai/ce-wifi-bt-test.service
[Unit]
Description=CE wifi/bt testing

[Service]
User=root
ExecStart=/bin/bash /usr/bin/ce-wifi-bt-test.sh

[Install]
WantedBy=multi-user.target
