# recipes-bsp/alsa-state/alsa-state/imx8mm-jaguar-dt510/tas6424-init.service
[Unit]
Description=DT510 TAS6424 class-D sane ALSA mixer defaults
After=sound.target
After=alsa-state.service
Wants=sound.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/tas6424-init
StandardOutput=journal
StandardError=journal
TimeoutStartSec=90
Restart=on-failure
RestartSec=3
StartLimitBurst=5
StartLimitInterval=120

[Install]
WantedBy=multi-user.target
