# recipes-bsp/alsa-state/alsa-state/imx8mm-jaguar-sentai/tas2563-init.service
[Unit]
Description=TAS2563 Audio Codec Initialization
After=sound.target alsa-state.service
Wants=sound.target
Requisite=sound.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/tas2563-init echo-removal
StandardOutput=journal
StandardError=journal
TimeoutStartSec=30
# Retry if audio card is not ready yet
Restart=on-failure
RestartSec=2
StartLimitBurst=5
StartLimitInterval=30

[Install]
WantedBy=multi-user.target
