# recipes-bsp/usb-composite-gadget/usb-composite-gadget/usb-composite-gadget.service
[Unit]
Description=USB Composite Gadget (CDC Serial + UAC2 Audio)
Documentation=man:systemd.service(5)
After=network.target
Wants=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/setup-usb-composite-gadget setup
ExecStartPost=/bin/bash -c 'echo device > /sys/devices/platform/soc@0/32c00000.bus/32e40000.usb/ci_hdrc.0/usb_role/ci_hdrc.0-role-switch/role'
ExecStop=/usr/bin/setup-usb-composite-gadget stop
User=root
StandardOutput=journal
StandardError=journal

# Restart policy
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
