# recipes-support/wifi-power-management/wifi-power-management/wifi-power-management.service
[Unit]
Description=WiFi Power Management for E-ink Board
After=network.target
Wants=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/wifi-power-management.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
