# recipes-bsp/udev/udev-rules-imx/imx8mm-jaguar-dt510/20-jaguar.rules
SUBSYSTEM=="tty", ATTRS{idVendor}=="0403", ATTRS{idProduct}=="6001", SYMLINK="ttyUSBConsole", RUN+="/sr/bin/stty -F /dev/ttyUSBConsole 115200"
