# recipes-bsp/eink-power-management/eink-power-management/eink-suspend.service
[Unit]
Description=E-ink Board Suspend Preparation
Before=sleep.target suspend.target
StopWhenUnneeded=yes

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/eink-suspend.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=30

[Install]
WantedBy=sleep.target suspend.target
