# recipes-bsp/eink-power-management/eink-power-management/eink-resume.service
[Unit]
Description=E-ink Board Resume Restoration
After=suspend.target sleep.target

[Service]
Type=oneshot
ExecStart=/usr/bin/eink-resume.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=30

[Install]
WantedBy=suspend.target sleep.target
