# recipes-bsp/eink-power-management/eink-power-management/eink-power-config.service
[Unit]
Description=E-ink Board Power Configuration
After=multi-user.target
Wants=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/eink-power-config.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
