# recipes-bsp/eink-power-management/eink-power-management/rtc-sync-time.path
[Unit]
Description=Monitor systemd-timesyncd for time synchronization events
Documentation=https://github.com/DynamicDevices/meta-dynamicdevices
After=systemd-timesyncd.service

[Path]
# Trigger when systemd-timesyncd state file changes (indicates time sync)
# This path is created/updated when timesyncd synchronizes time
PathChanged=/run/systemd/timesync/synchronized
# Also trigger on service state changes
Unit=systemd-timesyncd.service

[Install]
WantedBy=multi-user.target

