# recipes-bsp/eink-power-management/eink-power-management/wifi-connect.service
[Unit]
Description=Ensure WiFi Connection for E-Ink Board
Documentation=https://github.com/DynamicDevices/meta-dynamicdevices
# Wait for NetworkManager to be ready (NOT network-online.target - we're trying to establish connection!)
# network-online.target requires an active connection, which defeats the purpose
After=NetworkManager.service
Wants=NetworkManager.service
# Only start if WiFi interface exists
ConditionPathExists=/sys/class/net/wlan0

[Service]
Type=oneshot
# Wait for WiFi interface to be ready (NetworkManager may need time to initialize)
# Check every 3 seconds for faster response (30s timeout = 10 checks)
ExecStartPre=/bin/bash -c 'timeout=30; elapsed=0; while [ $elapsed -lt $timeout ]; do if [ -d /sys/class/net/wlan0 ] && nmcli -t -f STATE device show wlan0 2>/dev/null | grep -qE "disconnected|unavailable"; then break; fi; sleep 3; elapsed=$((elapsed+3)); done'
# Trigger immediate connection attempt for highest-priority WiFi connection (bypasses NetworkManager's long retry delay)
# NetworkManager's internal retry delay is ~60 seconds, which is too long for headless devices
# This ensures connection happens promptly on boot (within ~30 seconds)
# Find the highest priority WiFi connection with autoconnect enabled, or first WiFi connection if no priorities set
ExecStart=/bin/bash -c 'WIFI_CONN=$(nmcli -t -f NAME,TYPE,autoconnect,autoconnect-priority connection show | grep -E "wifi|802-11-wireless" | grep ":yes:" | sort -t: -k4 -rn | head -1 | cut -d: -f1); if [ -z "$WIFI_CONN" ]; then WIFI_CONN=$(nmcli -t -f NAME,TYPE connection show | grep -E "wifi|802-11-wireless" | head -1 | cut -d: -f1); fi; if [ -n "$WIFI_CONN" ]; then if ! nmcli -t -f STATE connection show "$WIFI_CONN" 2>/dev/null | grep -q "activated"; then echo "Triggering immediate WiFi connection to $WIFI_CONN (bypassing autoconnect delay)..."; nmcli connection up "$WIFI_CONN" || true; else echo "WiFi connection $WIFI_CONN already active"; fi; else echo "No WiFi connection configured"; fi'
StandardOutput=journal
StandardError=journal
TimeoutStartSec=60
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target

