# recipes-kernel/kernel-modules/kernel-module-nxp-wlan/imx8mm-jaguar-sentai/enable-wifi.service
[Unit]
Description=Enable WiFi GPIO

[Service]
User=root
ExecStart=/bin/bash /usr/bin/enable-wifi.sh

[Install]
WantedBy=multi-user.target
